# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
[program:x]
command=php /src/src/scripts/consumer.php
autostart=true
autorestart=true
startsecs=5
[program:phpfpm]
command=php-fpm
autostart=true
autorestart=false
